# etc/udev/rules.d/95-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/systemctl start monitor-connected.service"
